# rootfs-overlay/etc/inittab
# Respawn the Ouroboros daemons if they die
::respawn:/usr/bin/ouroboros_watchdog
::respawn:/usr/bin/ouroboros_remote
